# rootfs/etc/services.d/nginx/run
#!/usr/bin/execlineb -P
/opt/router/sbin/nginx
